# crontab/ana@dragocolosus.dyndns.org
# Mirror 
0 4 * * *  /home/ana/bin/photos/syncPhotos.pl >& /export/photos/log/cpAlbum.log
