# crontab/paus@ppc.mit.edu
# backup wordpress
28 1 * * * /home/paus/bin/save-it.py >& /home/paus/save-it.log
