# crontab/paus@t3desk004.mit.edu
# update repository mirros in redmine cache
*/5 * * * *     /home/paus/work/redmine/updateAll.sh >& /home/paus/work/redmine/lastUpdateAll.log
