# crontab/root@t3desk004.mit.edu
SHELL = /bin/tcsh
# Complete backup on Sundays, 3:00 am
0 3 * * 0       /backup/csh/run-complete-backup-to-logfile.csh    >> /dev/null
# Incremental any other day at 3:00 am
0 3 * * 1-6     /backup/csh/run-incremental-backup-to-logfile.csh >> /dev/null
# perform full backup of redmine database and attached files
15 04 * * *     /home/paus/bin/backup/backup-redmine.sh           >& /home/system/backup/redmine.log
# perform full backup of redmine database and attached files
20 04 * * *     /home/paus/bin/backup/backup-mysql.sh             >& /home/system/backup/mysql.log
