# crontab/root@t3btch039.mit.edu
# perform full backup of mysql database and attached files
20 04 * * * /home/paus/bin/backup/backup-mysql.sh /home/sysadmin/backup >& /home/system/backup/mysql.log
