# crontab/root@paushaus.dyndns.org
# Mirror user home directory to home server
0 5 * * *  /home/paus/bin/backup/syncHome.pl >> /dev/null
